# base-ubuntu-20.04/etc/service/socklog/run
#!/bin/sh
set -eu
. /usr/lib/runit-lib

exec s6-envuidgid nobody socklog unix /dev/log
